# Docs/kanxue-crawler/deploy/kanxue-crawler.service
[Unit]
Description=Kanxue forum crawler (one-shot, scheduled by .timer)
Documentation=file:///opt/kanxue-crawler/README.md
After=network-online.target
Wants=network-online.target
# 失败时发邮件由 Python 端在 scheduled_run 里直接 send_email，不走 systemd OnFailure。

[Service]
Type=oneshot
User=kanxue
Group=kanxue
WorkingDirectory=/opt/kanxue-crawler

EnvironmentFile=-/etc/kanxue-crawler.env

ExecStart=/opt/kanxue-crawler/.venv/bin/python -m src.scheduled_run --list-pages 1 --threads 100

TimeoutStartSec=4h

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/kanxue-crawler/output /opt/kanxue-crawler/state

StandardOutput=journal
StandardError=journal
SyslogIdentifier=kanxue-crawler

[Install]
WantedBy=multi-user.target
